# music_assistant/rootfs/etc/services.d/music_assistant/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: MusicAssistant
# Runs MusicAssistant
# ==============================================================================

# Enable Jemalloc
export LD_PRELOAD="/usr/local/lib/libjemalloc.so.2"

args="--config /config"

if bashio::config.true 'debug_messages'; then
    args="${args} --debug"
fi

exec python -m music_assistant ${args}